# conf/systemd/default.service
[Unit] # Unit部分，定义了单元的元数据和依赖关系
Description=Lightweight Music Server
After=network.target

[Service] # 服务(Service)部分，定义了服务的运行方式
Type=simple
Restart=on-failure # 当服务失败时自动重启服务
RestartSec=1
WorkingDirectory=/var/lms
ExecStart=/usr/bin/lms
User=lms
Group=lms
Environment=OMP_THREAD_LIMIT=1 # 设置环境变量，限制OpenMP线程数为1

[Install] # Install部分，定义了服务的安装信息
WantedBy=multi-user.target

